# scripts/launchagents/com.ystar.gov_mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ystar.gov_mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/Cellar/python@3.11/3.11.14_3/Frameworks/Python.framework/Versions/3.11/Resources/Python.app/Contents/MacOS/Python</string>
        <string>-m</string>
        <string>gov_mcp</string>
        <string>--session-config</string>
        <string>/Users/haotianliu/.openclaw/workspace/ystar-company/.ystar_session.json</string>
        <string>--transport</string>
        <string>sse</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>7922</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/haotianliu/.openclaw/workspace/ystar-company</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/haotianliu/.openclaw/workspace/Y-star-gov</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/ystar_gov_mcp.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ystar_gov_mcp_error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
