# scripts/k9_daily_patrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ystar.k9_daily_patrol</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/haotianliu/.openclaw/workspace/ystar-company/scripts/k9_daily_patrol.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/k9_daily_patrol_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/k9_daily_patrol_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/haotianliu/.openclaw/workspace/ystar-company</string>

    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
